# systemd/skype.service
[Unit]
Description=skype
Before=progs.target

[Service]
ExecStart=/usr/bin/skype
Restart=no

[Install]
WantedBy=progs.target
